# systemd/.config/systemd/user/bat.timer
[Unit]
Description=battery check timer

[Timer]
OnBootSec=120sec
OnUnitActiveSec=120sec
AccuracySec=5sec

[Install]
WantedBy=timers.target
