# systemd/.config/systemd/user/bat.timer
[Unit]
Description=probe battery level

[Timer]
OnUniActiveSec=45sec

[Install]
WantedBy=timers.target
