# systemd/.config/systemd/user/bat.timer
[Unit]
Description=battery check timer

[Timer]
OnBootSec=15sec
OnUnitActiveSec=15sec
AccuracySec=5sec

[Install]
WantedBy=timers.target
