# systemd/.config/systemd/user/bat.service
[Unit]
Description=battery check daemon

[Service]
ExecStart=/home/om/.config/systemd/user/bat.sh
