# systemd/.config/systemd/user/bat.timer
[Unit]
Description=battery check timer

[Timer]
OnBootSec=30sec
OnUnitActiveSec=30sec
AccuracySec=2sec

[Install]
WantedBy=timers.target
