# systemd/.config/systemd/user/bat.service
[Unit]
Description=battery check daemon

[Service]
Type=simple
ExecStart=/home/om/scripts/bat_not.sh

[Install]
WantedBy=multi.unit.target
