# systemd/.config/systemd/user/clean_homedir.service
[Unit]
Description=delete Downloads, Desktop folders from home directory

[Service]
Type=oneshot
ExecStart=$HOME/.config/systemd/user/clean_homedir.sh

[Install]
WantedBy=multi.unit.target
